# services/boot-display.service
[Unit]
Description=Boot Display Service
DefaultDependencies=no
After=systemd-remount-fs.service
Before=sysinit.target
Conflicts=shutdown.target

[Service]
Type=simple
User=root
Environment=PYTHONPATH=/home/admin/Radiowecker
WorkingDirectory=/home/admin/Radiowecker
ExecStart=/usr/bin/python3 /home/admin/Radiowecker/boot_display.py
RemainAfterExit=no

[Install]
WantedBy=sysinit.target
