# services/boot-display.service
[Unit]
Description=Boot Display Service
DefaultDependencies=no
After=systemd-remount-fs.service
Before=sysinit.target
Conflicts=shutdown.target

[Service]
Type=simple
User=%USER%
Environment=PYTHONPATH=%HOME%/Radiowecker
WorkingDirectory=%HOME%/Radiowecker
ExecStart=/usr/bin/python3 boot_display.py
StandardOutput=journal
StandardError=journal
RemainAfterExit=no

[Install]
WantedBy=sysinit.target
