# services/boot-display.service
[Unit]
Description=Boot Display Service
DefaultDependencies=no
After=systemd-udevd.service
Before=basic.target
Conflicts=shutdown.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%HOME%/Radiowecker
ExecStartPre=/bin/sleep 1
ExecStart=/usr/bin/python3 %HOME%/Radiowecker/boot_display.py
RemainAfterExit=no

[Install]
WantedBy=basic.target

